# crates/kit/src/units/bcvk-var-ephemeral.service
[Unit]
Description=Setup ephemeral /var from image content
DefaultDependencies=no
ConditionPathExists=/etc/initrd-release
Before=initrd-fs.target
After=sysroot.mount initrd-parse-etc.service
Requires=sysroot.mount

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=60
# Use overlayfs like /etc to avoid copying large files
# from virtiofs to tmpfs. The cp -a approach fails on virtiofs due to xattr issues
# and is prohibitively slow for large /var trees.
# Bind-mount /sysroot/var to a separate location first (same technique as /etc overlay).
ExecStart=/usr/bin/mkdir -p /run/var-lower /run/var-upper /run/var-work
ExecStart=/usr/bin/mount --bind /sysroot/var /run/var-lower
ExecStart=/usr/bin/mount -t overlay overlay -o lowerdir=/run/var-lower,upperdir=/run/var-upper,workdir=/run/var-work,index=off,metacopy=off /sysroot/var
